# pitchey-monitor.service
[Unit]
Description=Pitchey Production Health Monitor
After=network.target

[Service]
Type=simple
User=pitchey
WorkingDirectory=/home/supremeisbeing/pitcheymovie/pitchey_v0.2
ExecStart=/usr/bin/node health-check-worker.js
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
Environment="NODE_ENV=production"

[Install]
WantedBy=multi-user.target
