# monitoring/performance/pitchey-dashboard-generator.service
[Unit]
Description=Pitchey Dashboard Generator
After=network.target
Wants=network.target

[Service]
Type=simple
User=supremeisbeing
WorkingDirectory=/home/supremeisbeing/pitcheymovie/pitchey_v0.2/monitoring/performance
Environment=API_URL=https://pitchey-production.ndlovucavelle.workers.dev
Environment=DASHBOARD_REFRESH_INTERVAL=1
ExecStart=/usr/bin/deno run --allow-net --allow-read --allow-write real-time-dashboard.ts --continuous
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
