# monitoring/performance/pitchey-health-monitor.service
[Unit]
Description=Pitchey Health Monitor Daemon
After=network.target
Wants=network.target

[Service]
Type=simple
User=supremeisbeing
WorkingDirectory=/home/supremeisbeing/pitcheymovie/pitchey_v0.2/monitoring/performance
Environment=API_URL=https://pitchey-production.cavelltheleaddev.workers.dev
Environment=HEALTH_CHECK_INTERVAL=5
Environment=WEBHOOK_URL=
ExecStart=/usr/bin/deno run --allow-net --allow-read --allow-write health-check-daemon.ts
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
